# debian/preinst
#!/bin/bash

if [ -f /etc/apt/preferences.d/00-systemd-priority ]; then
   rm /etc/apt/preferences.d/00-systemd-priority;
fi

systemctl disable powertop
#systemctl disable powertop-usb-mouse
